# templates/rc.local.erb
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

get_iana_ip_address ()
{
  local NETBLOCK
  local IP_ADDRESS

  # If there's no default route (possible with host-only networking), then
  # look for devices assigned to an IANA private network with preference
  # given to one assigned a 192.168.x.y address.
  for NETBLOCK in 192.168.0.0/16 172.16.0.0/24 10.0.0.0/8; do
    IP_ADDRESS=$(ip route show to root $NETBLOCK | awk 'NR == 1 {print $9}')
    [ ! -z $IP_ADDRESS ] && echo $IP_ADDRESS && exit
  done

  # Nothing yet
  exit
}

get_primary_interface ()
{
  echo $(ip route show default scope global | awk 'NR==1 {print $5}')
}

<% if @hostname == "learning" %>
get_ip ()
{
  local IP_ADDRESS
  local PRIMARY_INTERFACE=$(get_primary_interface)

  if [ ! -z "$PRIMARY_INTERFACE" ]; then
    IP_ADDRESS=$(ip -o -4 a show $PRIMARY_INTERFACE | egrep '(en|eth)' | awk 'BEGIN {FS = "[ /]+"} ; {print $4}') || exit
    echo $IP_ADDRESS
  else
    IP_ADDRESS=$(get_iana_ip_address)
  fi
}

#Wait for network before loading IP address
COUNTER=0
until
  get_ip | grep '[0-9]' || [ $COUNTER -gt 12 ]
do
  sleep 5
  let COUNTER+=1
done

IP_ADDRESS=$(get_ip)
<% else %>
PRIMARY_INTERFACE=$(get_primary_interface)
if [ ! -z "$PRIMARY_INTERFACE" ]; then
  IP_ADDRESS=$(ip -o -4 a show $PRIMARY_INTERFACE | awk 'BEGIN {FS = "[ /]+"} ; {print $4}')
else
  IP_ADDRESS=$(get_iana_ip_address)
fi
<% end %>

if [ ! -s /var/local/password ]
then

  # Generate password in format WORD.WORD
  echo $(shuf -n 1 --random-source=/dev/urandom /var/local/places.txt).$(shuf -n 1 --random-source=/dev/urandom /var/local/places.txt) > /var/local/password
  echo "root:$(cat /var/local/password)" | chpasswd
fi

PASSWORD=$(cat /var/local/password)


echo \
"                --- Welcome to the PuppetLabs <%= @hostname.capitalize -%> VM ---

<%= @login_message %>

Use SSH to log in with the randomly generated password below.

    IP: $IP_ADDRESS
    username: root
    password: $PASSWORD

The password can be found in /var/local/password as well as this screen.
To set a new password, run the command: reset_password.rb
" > /etc/issue

# Disable the floppy driver
rmmod floppy

/root/.ssh_keygen.sh &
