# templates/rc.local.erb
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

PRIMARY_INTERFACE=$(ip route show default scope global | awk '{print $5}')
IP_ADDRESS=$(ip -o -4 a show $PRIMARY_INTERFACE | awk '{print $4}' | cut -d'/' -f1)

if [ ! -f /var/lock/password ]
then
<% if @operatingsystemmajrelease == "7" %>
  GPL_PATH=/usr/share/doc/centos-release/
<% else %>
  GPL_PATH=/usr/share/doc/centos-release-6/
<% end %>

  # Generate wordlist from GPL
  for word in $(shuf -n 1 $GPL_PATH/GPL)
  do
    echo $word | grep -v '^\s$' >> /tmp/wordlist
  done

  # Generate password in format WORD.Number.WORD
  echo $(shuf -n 1 /tmp/wordlist).$(echo $RANDOM | cut -c-2).$(shuf -n 1 /tmp/wordlist) > /var/lock/password
  echo "root:$(cat /var/lock/password)" | chpasswd
fi

PASSWORD=$(cat /var/lock/password)


echo \
"                --- Welcome to the PuppetLabs <%= @hostname.capitalize -%> VM ---

<%= @login_message %>






For the best experience please log in via SSH

    IP: $IP_ADDRESS
    username: root
    password: $PASSWORD

Your password was randomly generated from the text of the GNU public license
" > /etc/issue

/root/.ssh_keygen.sh &
