# templates/rc.local.erb
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

PRIMARY_INTERFACE=$(ip route show default scope global | awk '{print $5}')
IP_ADDRESS=$(ip -o -4 a show eth0 | awk '{print $4}' | cut -d'/' -f1)

echo \
"                --- Welcome to the PuppetLabs <%= @hostname.capitalize -%> VM ---

<%= @login_message %>






 For the best experience please log in via SSH

    IP: $IP_ADDRESS
    username: root
    password: puppet

" > /etc/issue

/root/.ssh_keygen.sh &

